# docker/supervisord/supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
priority = 1
command = /usr/sbin/sshd -D

[program:influxdb]
priority = 2
command = /usr/bin/influxd -pidfile /var/run/influxdb/influxd.pid -config /etc/influxdb/influxdb.conf

[program:syslog]
priority = 3
command = bash -c "sleep 30 && /etc/init.d/syslog-ng start"

[program:telegraf]
priority = 4
command = bash -c "sleep 30 && /etc/init.d/telegraf start"

[program:grafana]
priority = 5
command = /etc/init.d/grafana-server start

[program:chronograf]
priority = 6
command = /usr/bin/chronograf -r

[program:snmpwalk]
priority = 7
command = snmpwalk -c public -v 2c 172.17.0.1

